# conf/supervisord.conf
[unix_http_server]
file = /tmp/supervisor.sock

[supervisord]
logfile=/var/log/supervisord/supervisord.log
childlogdir=/var/log/supervisord/
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=%(run_daemons_as_user)s

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:%(site_name)s]
directory=%(project_dir)s
command=/usr/local/bin/gunicorn_django
autostart=true
autorestart=true
redirect_stderr=True